# config/systemd/user/at_notify.service
[Unit]
Description=Notify on at job low

[Service]
ExecStart=/home/jonathan/.local/bin/at_notify.sh


[Install]
WantedBy=default.target
