# config/systemd/user/at_notify.path
[Unit]
Description=Watch actual file for at notification
PartOf=at_notify.service

[Path]
PathModified=/tmp/remind

[Install]
WantedBy=default.target
