# config/systemd/user/battery_notify.path
[Unit]
Description=Watch actual file for battery notification
PartOf=battery_notify.service

[Path]
PathModified=/tmp/bat

[Install]
WantedBy=default.target
