# config/systemd/user/downloads.path
[Unit]
Description=Notify on download
PartOf=downloads.service

[Path]
PathModified=/home/jonathan/Downloads

[Install]
WantedBy=default.target
