# config/systemd/user/downloads.service
[Unit]
Description=Trying things

[Service]
ExecStart=/usr/bin/notify-send 'Download completed'

[Install]
WantedBy=default.target
