# config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/jonathan/.nix-profile/bin/mbsync -Va
