# config/systemd/user/mbsync.path
[Unit]
Description=Notify on download
PartOf=mbsync.service

[Path]
PathModified=/home/jonathan/.config/mail

[Install]
WantedBy=default.target
