# config/systemd/user/mbsync_clean.path
[Unit]
Description=Clean up weird artifacts
PartOf=mbsync_clean.service

[Path]
PathModified=/home/jonathan/.config/mail

[Install]
WantedBy=default.target
