# config/systemd/user/mu.service
[Unit]
Description=Mailbox indexer

[Service]
Type=oneshot
ExecStart=/home/jonathan/.nix-profile/bin/mu index
