# config/systemd/user/mu.service
[Unit]
Description=Mailbox indexer

[Service]
Type=oneshot
ExecStart=/run/current-system/sw/bin/mu index
