# config/systemd/user/mu.service
[Unit]
Description=Mailbox indexer

[Service]
Type=oneshot
ExecStart=/usr/bin/mu index
