# config/systemd/user/mu.timer
[Unit]
Description=Mailbox indexer timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=17m
Unit=mu.service

[Install]
WantedBy=timers.target
