# config/systemd/user/remind.service
[Unit]
Description=Reminder notifications

[Service]
Type=oneshot
ExecStart=/home/jonathan/.local/bin/remind.sh
