# config/systemd/user/update.service
[Unit]
Description=Update flakes daily

[Service]
Type=oneshot
ExecStart=/home/jonathan/.local/systemd/flake_update
