# config/systemd/user/update_notify.service
[Unit]
Description=Notify on battery low

[Service]
ExecStart=/home/jonathan/.local/systemd/update_notify


[Install]
WantedBy=default.target
