# config/systemd/user/update.timer
[Unit]
Description=Check for updates every 3 hours

[Timer]
OnBootSec=2m
OnUnitActiveSec=3h
Unit=update.service

[Install]
WantedBy=timers.target
